# poky/meta-qti-bsp/recipes-core/busybox/files/busybox-syslog.service
[Unit]
Description=System Logging Service

[Service]
Type=forking
RemainAfterExit=yes
ExecStart=/etc/init.d/syslog start
ExecStop=/etc/init.d/syslog stop

[Install]
WantedBy=multi-user.target
