# poky/meta-qti-data/recipes/data-ipa-cfg-mgr/files/ipacm.service
[Unit]
SourcePath=/usr/bin/ipacm
Description= IPACM Service

[Service]
Restart=no
ExecStart=/usr/bin/ipacm

[Install]
WantedBy=multi-user.target
