# poky/meta-vicos/recipes-extended/wiretestprog/wiretestprog/vic-test.service
[Unit]
Description=Vector test program

[Service]
Type=Simple
ExecStart=/bin/vic-test

[Install]
WantedBy=multi-user.target