# poky/meta-qti-bsp/recipes-connectivity/chrony/chrony/chronyd.path
[Unit]
Description=Monitor wifi connectivity to start chronyd

[Path]
PathExists=/tmp/notify-wifi-connected

[Install]
WantedBy=multi-user.target
